# ansible/roles/platform_services/templates/homelab-platform.service.j2
[Unit]
Description=Homelab Platform Services (OpenSearch + Dashboards + Vault)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ compose_dir }}
ExecStart=/usr/bin/podman-compose up -d
ExecStop=/usr/bin/podman-compose down
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
